# units/ecalserver.service.in
[Unit]
Description=Web server for the ESports calendar
After=network.target

[Service]
Type=simple
ExecStart=BIN_DIR/ecalserver --path DATA_DIR start --debug

[Install]
WantedBy=default.target
