# units/ecalserver.service.in
[Unit]
Description=Web server for the ESports calendar
After=network.target

[Service]
Type=simple
ExecStart=BIN_DIR/ecalserver start --debug --path DATA_DIR

[Install]
WantedBy=default.target
