# units/ecalserver.service.in
[Unit]
Description=Web server for the ESports calendar
After=network.target

[Service]
WorkingDirectory=DATA_DIR
Type=simple

ExecStart=BIN_DIR/ecalserver start --debug

