# units/events.service.in
[Unit]
Description=Scraper for ESports calendar
After=network.target

[Service]
Type=oneshot

ExecStart=BIN_DIR/ecalctl --path DATA_DIR fetch --debug --calendar CALENDARS

