# units/tooter.service.in
[Unit]
Description=Mastodon tooter for ESport events
After=network.target

[Service]
WorkingDirectory=DATA_DIR
Type=oneshot

ExecStart=BIN_DIR/BIN_CTL toot --type CALENDARS

